# scripts/release-assets/systemd/agent-container-hub.service
[Unit]
Description=agent-container-hub
After=network.target
Wants=network.target

[Service]
Type=simple
User=agent-container-hub
WorkingDirectory=/opt/agent-container-hub
EnvironmentFile=/opt/agent-container-hub/.env
ExecStart=/opt/agent-container-hub/agent-container-hub
Restart=on-failure
RestartSec=5
TimeoutStopSec=15
KillSignal=SIGTERM
LimitNOFILE=65535

[Install]
WantedBy=multi-user.target
